# systemd/voice-secretary-web.service
[Unit]
Description=Voice Secretary Web Dashboard
After=network.target

[Service]
Type=simple
User=voicesec
Group=voicesec
WorkingDirectory=/opt/voice-secretary
ExecStart=/opt/voice-secretary/.venv/bin/python -m app
Restart=on-failure
RestartSec=5
Environment=FLASK_ENV=production

[Install]
WantedBy=multi-user.target
