# systemd/voice-secretary-engine.service
[Unit]
Description=Voice Secretary SIP/AI Engine
After=network.target asterisk.service ollama.service
Wants=asterisk.service ollama.service

[Service]
Type=simple
User=voicesec
Group=voicesec
WorkingDirectory=/opt/voice-secretary
ExecStart=/opt/voice-secretary/.venv/bin/python -m engine
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
